# logger/.docker/service/push-to-redis/run
#!/bin/bash
/app/push-to-redis.php
sleep 1;
